# app/config/crontab
0 4 * * * /usr/src/app/scripts/import-prices.sh
